# deployment/containers/actual-budget/actual-budget.container
[Unit]
Description=Actual Budget - Personal Finance Manager
Wants=network.target
After=network.target

[Container]
ContainerName=actual-budget
Image=docker.io/actualbudget/actual-server:latest
AutoUpdate=registry

Network=web

PublishPort=5006:5006

# Mount volume for persistent data
Volume=%h/containers/actual-budget/data:/data

# Traefik labels for reverse proxy
Label=traefik.enable=true
# Dual-host: reachable via the public domain or the local network directly.
Label=traefik.http.routers.actual_budget.rule=Host(`budget.{{ oci_parent_host }}`) || Host(`budget.{{ local_host }}`)
Label=traefik.http.routers.actual_budget.entryPoints=websecure
Label=traefik.http.routers.actual_budget.tls=true
Label=traefik.http.routers.actual_budget.tls.certresolver=cloudflare
# Only request a cert for the public name - local_host has no public DNS, so
# an ACME challenge for it would just fail.
Label=traefik.http.routers.actual_budget.tls.domains[0].main=budget.{{ oci_parent_host }}
Label=traefik.http.services.actual_budget.loadbalancer.server.port=5006

[Service]
Restart=always

[Install]
WantedBy=default.target
