# deployment/containers/crowdsec-bouncer/crowdsec-bouncer.container
[Unit]
Description=CrowdSec Traefik Bouncer
Wants=network.target crowdsec.service
After=network.target crowdsec.service

[Container]
ContainerName=crowdsec-bouncer
Image=docker.io/fbonalair/traefik-crowdsec-bouncer:latest
AutoUpdate=registry

Network=web

Environment=CROWDSEC_AGENT_HOST=crowdsec:8080
Environment=CROWDSEC_BOUNCER_ORIGIN_HEADER=X-Forwarded-For
Environment=GIN_MODE=release

# API key must first be generated on the crowdsec container via:
#   podman exec crowdsec cscli bouncers add crowdsec-bouncer
# then stored as crowdsec_bouncer_api_key in secrets.yml
Secret=crowdsec_bouncer_api_key,type=env,target=CROWDSEC_BOUNCER_API_KEY

# Register a reusable forwardAuth middleware other routers opt into via
# Label=traefik.http.routers.<name>.middlewares=crowdsec-bouncer@docker
Label=traefik.enable=true
Label=traefik.http.middlewares.crowdsec-bouncer.forwardauth.address=http://crowdsec-bouncer:8080/api/v1/forwardAuth
Label=traefik.http.middlewares.crowdsec-bouncer.forwardauth.trustForwardHeader=true

[Service]
Restart=always

[Install]
WantedBy=default.target
